# target/linux/ramips/base-files/etc/init.d/wireless
#!/bin/sh /etc/rc.common
#
# bygomma :  /etc/init.d/wireless
START=21
STOP=89

USE_PROCD=1

start_service()
{
	echo ">>>>> start_service()"
	procd_open_instance
	procd_set_param stdout 1 # forward stdout of the command to logd
	procd_set_param stderr 1 # same for stderr
	procd_close_instance
}

reload_service()
{
	logger -t "init.d/wireless" "reload_service"
	ubus call network reload
}

stop()
{
	logger -t "init.d/wireless" "stop"
}

service_running()
{
	logger -t "init.d/wireless" "service_running"
	ubus -t 30 wait_for network.interface
}

service_triggers()
{
	logger -t "init.d/wireless" "service_triggers"
	procd_add_reload_interface_trigger "lan"
}

restart()
{
	logger -t "init.d/wireless" "restart"
	stop "$@"
	start "$@"
}

shutdown()
{
	logger -t "init.d/wireless" "shutdown"
	echo ">>>>> shutdown()"
}

# bygomma : start
